# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:migrate]
command=bash -c "cd /app/backend && sleep 5 && alembic upgrade head 2>&1 && echo 'Migration complete' || echo 'Migration failed (may already be applied)'"
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
startsecs=0

[program:backend]
command=bash -c "cd /app/backend && sleep 10 && uvicorn main:app --host 0.0.0.0 --port 8000 --workers 2"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
startsecs=5

[program:celery]
command=bash -c "cd /app/backend && celery -A tasks.celery_app worker --loglevel=info --concurrency=2"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
startsecs=10

[program:frontend]
command=bash -c "cd /app && HOSTNAME=0.0.0.0 PORT=3000 exec node server-custom.js"
environment=NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=5
